# ubuntu-nginx-php-2004/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/etc/init.d/php7.4-fpm start

[program:nginx]
command=nginx