# ubuntu-gcc-1804/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/etc/init.d/ssh start
